# gdnsd/dns.conf
[program:dns]
command=/massbit/massbitroute/gdnsd/run.sh
autorestart=true
[program:dns_monitor]
command=/massbit/massbitroute/gdnsd/run.sh _monitor
autorestart=true
